# templates/default/docker.socket.erb
[Unit]
Description=Docker Socket for the API
PartOf=docker.service

[Socket]
<% Array(node['docker']['host']).each do |socket| -%>
<% next if socket.match(/^unix:/) -%>
<% next if socket.match(/^tcp:/) -%>
<% next if socket.match(/^fd:/) -%>
ListenStream=<%= socket %>
<% end -%>
ListenStream=/var/run/docker.sock
SocketMode=0660
SocketUser=root
SocketGroup=<%= node['docker']['group'] || 'root' %>

[Install]
WantedBy=sockets.target
