# templates/default/docker.conf.erb
description     "Run dockerd"

start on runlevel [2345]
stop on runlevel [!2345]

# if you want it to automatically restart if it crashes, leave the next line in
respawn

script
    <% case node['platform'] -%>
    <% when 'fedora' -%>
    test -f /etc/locale.conf && . /etc/locale.conf || true
    <% when 'ubuntu' -%>
    test -f /etc/default/locale && . /etc/default/locale || true
    <% end -%>
    test -f <%= node['docker']['config_dir'] %>/docker && . <%= node['docker']['config_dir'] %>/docker || true
    LANG=$LANG LC_ALL=$LANG <%= node['docker']['install_dir'] %>/docker $DOCKER_OPTS
end script
