# templates/default/docker-container.conf.erb
description "Docker service for <%= @service_name %>"
author "Docker Chef Cookbook"

start on filesystem and started docker
stop on runlevel [!2345]

# We don't want to TERM the `docker wait` process so we fake the signal
# we send to it. The pre-stop script issues the `docker stop` command
# which causes the `docker wait` process to exit
kill signal CONT

# Due to a bug in upstart we need to set the modes we consider
# successful exists https://bugs.launchpad.net/upstart/+bug/568288
normal exit 0 CONT

respawn

script
  # Wait for docker to finish starting up first.
  <% # DEPRECATED: bind_socket usage -%>
  <% if node['docker']['bind_socket'] -%>
  FILE=<%= node['docker']['bind_socket'].gsub(/^unix:\/\//,'') %>
  while [ ! -e $FILE ] ; do
    inotifywait -t 2 -e create $(dirname $FILE)
  done
  <% elsif Array(node['docker']['host']).any?{|host| host.start_with? 'unix://'} -%>
  FILE=<%= Array(node['docker']['host']).select{|host| host.start_with? 'unix://'}[0].gsub(/^unix:\/\//,'') %>
  while [ ! -e $FILE ] ; do
    inotifywait -t 2 -e create $(dirname $FILE)
  done
  <% end -%>
  <%= Docker::Helpers.executable(node) %> start -a <%= @service_name %>
end script

pre-stop script
  <%= Docker::Helpers.executable(node) %> stop -t <%= @cmd_timeout %> <%= @service_name %>
end script
