# templates/default/docker-container.socket.erb
[Unit]
Description=Docker sockets for <%= @service_name %>

[Socket]
<% @sockets.each do |s| -%>
ListenStream=<%= s %>
<% end -%>

[Install]
WantedBy=sockets.target
