# templates/default/docker.service.erb
[Unit]
Description=Docker container management daemon

[Service]
Type=simple
ExecStartPre=/usr/sbin/sysctl -w net.ipv4.ip_forward=1 net.ipv6.conf.all.forwarding=1
<% case node['platform'] -%>
<% when 'fedora' -%>
ExecStartPre=test -f /etc/locale.conf && . /etc/locale.conf || true
<% when 'ubuntu' -%>
ExecStartPre=test -f /etc/default/locale && . /etc/default/locale || true
<% end -%>
ExecStartPre=test -f <%= node['docker']['config_dir'] %>/docker && . <%= node['docker']['config_dir'] %>/docker || true
ExecStart=LANG=$LANG LC_ALL=$LANG <%= node['docker']['install_dir'] %>/docker $DOCKER_OPTS
ExecStartPost=/usr/sbin/iptables -I FORWARD -o docker0 -m conntrack --ctstate RELATED,ESTABLISHED -j ACCEPT
ExecStartPost=/usr/sbin/iptables -I FORWARD -i docker0 ! -o docker0 -j ACCEPT
Restart=on-failure

[Install]
WantedBy=multi-user.target
